# python-backend/textfind.service
[Unit]
Description=TextFind Python OCR Backend
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/textfind/python-backend
Environment="PATH=/opt/textfind/venv/bin"
EnvironmentFile=/opt/textfind/python-backend/.env
ExecStart=/opt/textfind/venv/bin/gunicorn --workers 2 --bind 0.0.0.0:5000 --timeout 120 ocr_server:app
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
